# parse_test/shell.service
[Unit]
Description=Shell

[Service]
Type=simple
ExecStart=/bin/shell.el
Restart=always
RestartSec=5000