# parse_test/shell.service
[Unit]
Description=Shell
After=http_server.service

[Service]
Type=simple
ExecStart=/bin/shell.elf
Restart=always