# parse_test/shell.service
[Unit]
Description=Shell

[Service]
Type=simple
ExecStart=/bin/shell.elf
Restart=always
RestartSec=5000