# parse_test/test_timer_unit.service
[Unit]
Description=a test for timer_unit

[Service]
Type=simple
ExecStart=-/bin/about.elf
ExecStartPre=-/bin/clear