# parse_test/test_timer_unit.timer
[Unit]
Description=a timer for the same_name.service

[Timer]
OnActiveSec=3s
OnUnitActiveSec=5s
Unit=test_timer_unit.service